# 1_data_collection/api_calls/crontab
0 10 * * * python /app/adzuna_ads_yesterday.py >> /var/log/adzuna_cron.log 2>&1
