# 1_data_collection/api_calls/crontab
0 10 * * * python3 /app/adzuna_ads_yesterday.py >> /app/adzuna_cron.log 2>&1



