# 1_data_collection/api_calls/crontab
TZ=Europe/Berlin
0 10 * * * python /app/adzuna_ads_yesterday.py >> /app/adzuna_cron.log 2>&1



